# udev-rules/usb-trigger.rules
ACTION=="add", RUN+="../main.py"
# update RUN script target
